# etc/sv/runsvdir-igor/run
#!/bin/sh

export USER="igor"
export HOME="/home/igor"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/.config/runit/service"

exec chpst -u "$USER:$groups" runsvdir "$svdir"
